# puppet/sp_worker/templates/wso2sp-worker.service.erb
[Unit]
Description=WSO2 Stream Processor - Worker
After=network.target

[Service]
ExecStart=<%=@install_path%>/bin/<%=@product_profile%>.sh start
ExecStop=<%=@install_path%>/bin/<%=@product_profile%>.sh stop
ExecRestart=<%=@install_path%>/bin/<%=@product_profile%>.sh restart
PIDFile=<%=@install_path%>/wso2/<%=@product_profile%>/runtime.pid
User=<%=@user%>
Group=<%=@user_group%>
Type=forking
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
